# laravel/php/crontab
* * * * * cd /var/www && php artisan schedule:run >> /home/cesarobedfl/cron.txt 2>&1